# usr/lib/opkg/info/collectd.list
/usr/share/collectd/types.db
/etc/init.d/collectd
/etc/config/collectd
/usr/sbin/collectd
/lib/upgrade/keep.d/collectd
