# usr/lib/opkg/info/flow_broker.list
/etc/init.d/flow_broker
/usr/share/flow_broker/flow_broker.py
/usr/sbin/flow_broker.sh
/etc/config/flow_broker
