# usr/lib/opkg/info/luci-app-statistics.list
/www/luci-static/resources/view/statistics/plugins/irq.js
/usr/share/luci/statistics/plugins/rrdtool.json
/www/luci-static/resources/view/statistics/plugins/processes.js
/www/luci-static/resources/statistics/rrdtool/definitions/processes.js
/usr/share/luci/statistics/plugins/contextswitch.json
/usr/share/luci/statistics/plugins/dhcpleases.json
/www/luci-static/resources/statistics/rrdtool/definitions/ipstatistics.js
/www/luci-static/resources/view/statistics/plugins/contextswitch.js
/www/luci-static/resources/statistics/rrdtool/definitions/iptables.js
/www/luci-static/resources/view/statistics/plugins/memory.js
/usr/bin/stat-genconfig
/www/luci-static/resources/statistics/rrdtool/definitions/curl.js
/www/luci-static/resources/view/statistics/plugins/dhcpleases.js
/usr/share/luci/statistics/plugins/uptime.json
/usr/share/luci/statistics/plugins/cpufreq.json
/usr/share/luci/statistics/plugins/nut.json
/usr/share/luci/statistics/plugins/syslog.json
/www/luci-static/resources/statistics/rrdtool/definitions/irq.js
/usr/share/luci/statistics/plugins/memory.json
/usr/share/luci/statistics/plugins/exec.json
/etc/uci-defaults/40_luci-statistics
/usr/share/luci/statistics/plugins/interface.json
/www/luci-static/resources/statistics/rrdtool/definitions/sqm.js
/www/luci-static/resources/view/statistics/plugins/disk.js
/www/luci-static/resources/statistics/rrdtool/definitions/tcpconns.js
/usr/share/luci/statistics/plugins/openvpn.json
/www/luci-static/resources/statistics/rrdtool/definitions/iwinfo.js
/usr/share/luci/statistics/plugins/tcpconns.json
/www/luci-static/resources/view/statistics/plugins/netlink.js
/www/luci-static/resources/statistics/rrdtool/definitions/chrony.js
/www/luci-static/resources/view/statistics/plugins/exec.js
/www/luci-static/resources/view/statistics/plugins/cpu.js
/www/luci-static/resources/view/statistics/plugins/iptables.js
/www/luci-static/resources/statistics/rrdtool/definitions/dhcpleases.js
/www/luci-static/resources/statistics/rrdtool/definitions/apcups.js
/www/luci-static/resources/statistics/rrdtool/definitions/thermal.js
/usr/share/luci/statistics/plugins/network.json
/www/luci-static/resources/statistics/rrdtool/definitions/entropy.js
/usr/share/luci/statistics/plugins/processes.json
/www/luci-static/resources/statistics/rrdtool/definitions/cpu.js
/www/luci-static/resources/view/statistics/plugins/interface.js
/www/luci-static/resources/statistics/rrdtool/definitions/interface.js
/usr/share/luci/statistics/plugins/thermal.json
/usr/share/luci/statistics/plugins/load.json
/www/luci-static/resources/view/statistics/plugins/snmp6.js
/usr/share/luci/statistics/plugins/df.json
/usr/share/luci/statistics/plugins/chrony.json
/www/luci-static/resources/view/statistics/plugins/tcpconns.js
/usr/share/luci/statistics/plugins/conntrack.json
/www/luci-static/resources/view/statistics/plugins/uptime.js
/etc/config/luci_statistics
/www/luci-static/resources/view/statistics/plugins/conntrack.js
/www/luci-static/resources/statistics/rrdtool/definitions/conntrack.js
/www/luci-static/resources/view/statistics/plugins/unixsock.js
/www/luci-static/resources/view/statistics/plugins/email.js
/usr/share/luci/statistics/plugins/dns.json
/www/luci-static/resources/statistics/rrdtool/definitions/disk.js
/usr/share/luci/statistics/plugins/netlink.json
/www/luci-static/resources/view/statistics/plugins/ping.js
/usr/share/luci/menu.d/luci-app-statistics.json
/usr/share/luci/statistics/plugins/disk.json
/www/luci-static/resources/view/statistics/plugins/nut.js
/www/luci-static/resources/view/statistics/collectd.js
/www/luci-static/resources/statistics/rrdtool/definitions/cpufreq.js
/www/luci-static/resources/statistics/rrdtool/definitions/splash_leases.js
/www/luci-static/resources/statistics/rrdtool/definitions/sensors.js
/www/luci-static/resources/statistics/rrdtool.js
/usr/share/luci/statistics/plugins/sensors.json
/usr/share/luci/statistics/plugins/entropy.json
/www/luci-static/resources/statistics/rrdtool/definitions/nut.js
/www/luci-static/resources/view/statistics/plugins/syslog.js
/www/luci-static/resources/view/statistics/plugins/thermal.js
/www/luci-static/resources/view/statistics/plugins/entropy.js
/www/luci-static/resources/view/statistics/plugins/ipstatistics.js
/www/luci-static/resources/statistics/rrdtool/definitions/openvpn.js
/usr/share/luci/statistics/plugins/csv.json
/www/luci-static/resources/statistics/rrdtool/definitions/memory.js
/usr/share/luci/statistics/plugins/snmp6.json
/www/luci-static/resources/view/statistics/plugins/olsrd.js
/www/luci-static/resources/view/statistics/plugins/iwinfo.js
/www/luci-static/resources/statistics/rrdtool/definitions/ping.js
/www/luci-static/resources/view/statistics/plugins/chrony.js
/etc/init.d/luci_statistics
/www/luci-static/resources/view/statistics/plugins/csv.js
/www/luci-static/resources/statistics/rrdtool/definitions/snmp6.js
/www/luci-static/resources/statistics/rrdtool/definitions/contextswitch.js
/www/luci-static/resources/view/statistics/plugins/load.js
/www/luci-static/resources/view/statistics/plugins/network.js
/www/luci-static/resources/view/statistics/plugins/apcups.js
/www/luci-static/resources/view/statistics/plugins/rrdtool.js
/www/luci-static/resources/view/statistics/plugins/cpufreq.js
/www/luci-static/resources/view/statistics/plugins/sensors.js
/usr/share/luci/statistics/plugins/iptables.json
/usr/share/rpcd/acl.d/luci-app-statistics.json
/usr/share/luci/statistics/plugins/curl.json
/usr/share/luci/statistics/plugins/splash_leases.json
/www/luci-static/resources/view/statistics/plugins/curl.js
/www/luci-static/resources/view/statistics/plugins/splash_leases.js
/www/luci-static/resources/statistics/rrdtool/definitions/sqmcake.js
/www/luci-static/resources/statistics/rrdtool/definitions/df.js
/usr/share/luci/statistics/plugins/cpu.json
/usr/share/luci/statistics/plugins/iwinfo.json
/www/luci-static/resources/statistics/rrdtool/definitions/uptime.js
/www/luci-static/resources/statistics/rrdtool/definitions/netlink.js
/www/luci-static/resources/view/statistics/plugins/openvpn.js
/usr/share/luci/statistics/plugins/olsrd.json
/usr/share/luci/statistics/plugins/ping.json
/www/luci-static/resources/statistics/rrdtool/definitions/load.js
/www/luci-static/resources/view/statistics/plugins/dns.js
/usr/share/luci/statistics/plugins/ipstatistics.json
/usr/share/luci/statistics/plugins/email.json
/www/luci-static/resources/view/statistics/graphs.js
/usr/share/luci/statistics/plugins/apcups.json
/usr/share/luci/statistics/plugins/irq.json
/www/luci-static/resources/view/statistics/plugins/df.js
/usr/share/luci/statistics/plugins/unixsock.json
/www/luci-static/resources/statistics/rrdtool/definitions/olsrd.js
/www/luci-static/resources/statistics/rrdtool/definitions/dns.js
/www/luci-static/resources/statistics/rrdtool/definitions/ip6tables.js
