# usr/lib/opkg/info/mwan3.list
/etc/mwan3.user
/etc/hotplug.d/iface/15-mwan3
/usr/sbin/mwan3
/lib/mwan3/common.sh
/etc/init.d/mwan3
/etc/hotplug.d/iface/16-mwan3-user
/etc/uci-defaults/mwan3-migrate-flush_conntrack
/usr/libexec/rpcd/mwan3
/usr/sbin/mwan3track
/lib/mwan3/mwan3.sh
/lib/mwan3/libwrap_mwan3_sockopt.so.1.0
/usr/sbin/mwan3rtmon
/etc/config/mwan3
