# usr/lib/opkg/info/netify-fwa.list
/etc/netify-fwa/netify-categories.json
/usr/share/netify-fwa/nfa_daemonize.py
/usr/share/netify-fwa/nfa_task.py
/usr/share/netify-fwa/nfa_config.py
/usr/share/netify-fwa/nfa_util.py
/usr/share/netify-fwa/nfa_defaults.py
/usr/share/netify-fwa/nfa_timer.py
/usr/share/netify-fwa/nfa_netify_api.py
/usr/share/netify-fwa/nfa_fw_openwrt.py
/usr/share/netify-fwa/inotify_simple.py
/etc/init.d/netify-fwa
/usr/share/netify-fwa/nfa_ipset.py
/etc/netify-fwa/netify-fwa.ini
/usr/share/netify-fwa/nfa_version.py
/usr/share/netify-fwa/nfa_netifyd.py
/usr/share/netify-fwa/nfa_rule.py
/etc/netify-fwa/netify-fwa.json
/usr/share/netify-fwa/nfa_global.py
/usr/share/netify-fwa/nfa_fw_iptables.py
/etc/netify-fwa/app-proto-data.json
/usr/share/netify-fwa/inotify_test.py
/usr/sbin/netify-fwa
/usr/share/netify-fwa/nfa_main.py
