# usr/lib/opkg/info/openwisp-config.list
/etc/init.d/openwisp_config
/usr/sbin/openwisp-reload-config
/usr/lib/lua/openwisp/utils.lua
/usr/sbin/openwisp_config
/usr/sbin/openwisp-get-address
/etc/config/openwisp
/usr/sbin/openwisp-store-unmanaged
/usr/sbin/openwisp-remove-default-wifi
/usr/sbin/openwisp-update-config
/etc/openwisp/VERSION
/usr/sbin/openwisp-restore-unmanaged
/usr/lib/lua/openwisp/net.lua
/usr/sbin/openwisp-uci-autoname
