# usr/lib/opkg/info/openwisp-monitoring.list
/usr/sbin/openwisp-monitoring
/etc/config/openwisp-monitoring
/etc/init.d/openwisp-monitoring
