# usr/lib/opkg/info/ulogd.list
/usr/lib/ulogd/ulogd_raw2packet_BASE.so
/etc/ulogd.conf
/usr/sbin/ulogd
/etc/init.d/ulogd
